# deploy/wa-hub.service
[Unit]
Description=wa-hub-demo — self-hosted WhatsApp HTTP API (Baileys)
Documentation=https://github.com/Noam13-w/wa-hub-demo
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wahub
Group=wahub
WorkingDirectory=/srv/wa-hub-demo
EnvironmentFile=/srv/wa-hub-demo/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /srv/wa-hub-demo/src/index.js
Restart=always
RestartSec=5
# Match src/index.js shutdown timeout (12 s) with a small grace buffer so systemd
# doesn't SIGKILL us mid-shutdown. Webhook drain + WS close + Baileys end all need a moment.
TimeoutStopSec=15
SyslogIdentifier=wa-hub-demo
StandardOutput=journal
StandardError=journal

# ─── Hardening ────────────────────────────────────────────────────────
NoNewPrivileges=true
# Default new files/dirs to owner-only (0600/0700) — the WhatsApp session under
# data/auth/ is account credentials and must never be group/world-readable.
UMask=0077
ProtectSystem=strict
ProtectHome=tmpfs
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
SystemCallArchitectures=native
# SystemCallFilter is intentionally left UNSET by default. Node 20+ uses syscalls
# (io_uring, clone3) that systemd's @system-service allow-list blocks, which makes
# the process core-dump with status=31/SYS on boot. Leaving the syscall filter open
# keeps all install routes reliable; every other sandbox directive above
# (NoNewPrivileges, ProtectSystem=strict, empty CapabilityBoundingSet, …) stays in force.
# To opt into a seccomp filter SAFELY (blocked calls return EPERM instead of
# SIGSYS-killing the process), install deploy/wa-hub.hardening.conf.example as a
# drop-in and verify the service still starts — see that file for instructions.
ReadWritePaths=/srv/wa-hub-demo/data

# Resource caps — adjust if you expect heavy traffic.
MemoryMax=512M
TasksMax=256
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
